# telldus-core/tdadmin/05-tellstick.rules
ID_VENDOR_ID=="1781", SUBSYSTEM=="usb"  RUN+="${CMAKE_INSTALL_PREFIX}/share/telldus-core/helpers/udev.sh"

