# telldus-core/service/05-tellstick.rules
BUS=="usb", ID_VENDOR_ID=="1781" RUN+="/usr/share/telldus-core/helpers/udev.sh"

